# .config/systemd/user/notes-sync.timer
[Unit]
Description=Periodic notes git sync — every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=notes-sync.service

[Install]
WantedBy=timers.target
