# .config/systemd/user/notes-sync.service
[Unit]
Description=Periodic notes git sync (safety net)

[Service]
Type=oneshot
ExecStart=/home/luke/scripts/notes-sync
